# ops/launchd/com.marketcow.production.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.marketcow.production</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>/Users/androidjk/Library/Application Support/MarketCow/start-production.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MARKETCOW_PROFILE</key>
        <string>v2-production</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/androidjk/Library/Logs/MarketCow/production.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/androidjk/Library/Logs/MarketCow/production.error.log</string>
</dict>
</plist>
